# .config/autostart/ssh-add.desktop
[Desktop Entry]
Exec=/home/thomas/.config/autostart/ssh-add.sh"
Icon=dialog-scripts
Name=ssh-add
Path=
Type=Application
X-KDE-AutostartScript=true
